# sonarr/rootfs/services.d/sonarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Add-on: Sonarr
# Runs the Sonarr server
# ==============================================================================

bashio::log.info "Starting the Sonarr addon...."

cd /opt/NzbDrone || bashio::exit.nok "setup gone wrong!"

exec mono --debug NzbDrone.exe -nobrowser -data=/config
